# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

; ── FastAPI ──────────────────────────────────────────────────────────────
[program:api]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 1 --no-access-log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
priority=10

; Celery worker removed — BRPOP polling consumed ~518k Redis requests/month,
; exceeding the Upstash free-tier limit (500k/month).
; Tasks now run directly in-process via FastAPI BackgroundTasks.
; Internal endpoints (POST /api/internal/tasks/*) are still triggered by cron-job.org.
; Celery Beat was already removed previously for the same reason.
